# nimbus/confs/nimbus.cron
#!/bin/sh

su - nimbus -c '/usr/bin/pg_dump nimbus | bzip2 -9 > /var/nimbus/nimbus-sql.bz2';
su - nimbus -c 'PGPASSWORD=bacula /usr/bin/pg_dump bacula -U bacula -h 127.0.0.1 | bzip2 -9 > /var/nimbus/bacula-sql.bz2';



su - nimbus -c '/var/www/nimbus --update-graphs-data';
su - nimbus -c '/var/www/nimbus --delete-volumes';

