# apps/unix/client/debian/prerm
#!/bin/bash

/usr/sbin/update-rc.d -f nimbusclient remove;
